# debian/postinst
#!/bin/sh

update-rc.d wb-mqtt-sht1x defaults
invoke-rc.d wb-mqtt-sht1x start
